# deploy/hf-space/supervisord.conf
; supervisord runs the whole stack inside one HF Spaces container (as UID 1000).
; All logs go to the container stdout/stderr so they show in the HF "Logs" tab.
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=appuser

; Control interface so `supervisorctl -c /home/appuser/supervisord.conf status`
; works inside the container (view/restart individual processes).
[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ---- datastores (start first) --------------------------------------------
[program:postgres]
command=postgres -D /home/appuser/pgdata -p 5432
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis]
command=redis-server --port 6379 --bind 127.0.0.1 --save "" --appendonly no
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:chroma]
command=chroma run --host 127.0.0.1 --port 8001 --path /home/appuser/chroma --log-path /tmp/chroma.log
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---- API: wait for Postgres, migrate + seed, then serve -------------------
[program:api]
directory=/home/appuser/app_src/backend
command=sh -c 'until pg_isready -h 127.0.0.1 -p 5432 -U compliance >/dev/null 2>&1; do echo "waiting for postgres..."; sleep 1; done; alembic upgrade head; python scripts/seed.py || true; exec uvicorn app.main:app --host 127.0.0.1 --port 8000'
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---- Celery worker: wait for Postgres + Redis, then run -------------------
[program:worker]
directory=/home/appuser/app_src/backend
command=sh -c 'until pg_isready -h 127.0.0.1 -p 5432 -U compliance >/dev/null 2>&1; do sleep 1; done; exec celery -A app.worker.celery_app worker --loglevel=INFO'
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ---- Nginx: public entrypoint on 7860 -------------------------------------
[program:nginx]
command=nginx -c /home/appuser/nginx.conf
priority=40
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
